# sponsorblockcast.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/sponsorblockcast
DynamicUser=yes

[Install]
WantedBy=multi-user.target
