# sponsorblockcast.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/sponsorblockcast

[Install]
WantedBy=multi-user.target
